# systemd/display-wiki.service.template
[Unit]
Description=Display Wiki Documentation Site
After=network.target

[Service]
Type=simple
WorkingDirectory=PROJECT_DIR
ExecStart=/usr/bin/python3 -m http.server 8080 --bind 0.0.0.0 --directory PROJECT_DIR/site
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
